# crawling/crontab
0 * * * * python3 /crawling/src/main.py >> /crawling/logs/crontab.log 2>&1